# setup/jenkins.conf
[program:jenkins]
user = nobody
command = /home/nobody/jenkins.sh

[program:openssh]
user = root
command = /root/openssh.sh